# deploy/portfolio.service
[Unit]
Description=Portfolio Management System
After=syslog.target network.target mariadb.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu

# 작업 디렉토리
WorkingDirectory=/opt/portfolio

# JAR 파일 실행
ExecStart=/usr/bin/java \
    -Xms512m \
    -Xmx1024m \
    -Dspring.profiles.active=prod \
    -Dserver.port=8080 \
    -jar /opt/portfolio/portfolio-backend-1.0.0.jar

# 환경 변수 (실제 값은 /etc/systemd/system/portfolio.service.d/override.conf에 설정)
Environment="ADMIN_USERNAME=${ADMIN_USERNAME}"
Environment="ADMIN_PASSWORD=${ADMIN_PASSWORD}"
Environment="DB_PASSWORD=${DB_PASSWORD}"

# 로그 설정
StandardOutput=journal
StandardError=journal
SyslogIdentifier=portfolio

# 재시작 정책
Restart=always
RestartSec=10

# 보안 설정
NoNewPrivileges=true
PrivateTmp=true

# 리소스 제한
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
